# systemd/freezer-display.service
[Unit]
Description=Freezer Monitor — LCD Display Service
After=network.target
RequiresMountsFor=/data

[Service]
Type=simple
User=pi
ExecStart=/usr/bin/python3 -u /opt/freezerpi/display_service.py
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=5
ReadWritePaths=/run /data
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes

[Install]
WantedBy=multi-user.target
