# systemd/icebox-netwatchdog.service
[Unit]
Description=IceboxHero — Network Watchdog
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/iceboxhero/network_watchdog.sh
StandardOutput=journal
StandardError=journal
