# systemd/icebox-netwatchdog.timer
[Unit]
Description=IceboxHero — Network Watchdog Timer
Requires=icebox-netwatchdog.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
